# deploy/systemd/mailngine-api.service
[Unit]
Description=Mailngine API Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=mailngine
Group=mailngine
WorkingDirectory=/opt/mailngine
ExecStart=/opt/mailngine/server
EnvironmentFile=/opt/mailngine/.env
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mailngine
PrivateTmp=true

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
