# deploy/production/mailngine-api.service
[Unit]
Description=Mailngine API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mailngine
EnvironmentFile=/opt/mailngine/shared/.env
ExecStart=/opt/mailngine/current/server
ExecStartPost=/bin/sh -c 'for i in 1 2 3 4 5 6; do sleep 2 && curl -sf http://localhost:8090/health && exit 0; done; exit 1'
Restart=always
RestartSec=5
LimitNOFILE=65535

# Graceful shutdown: Go server drains connections for up to 30s
TimeoutStopSec=35
KillMode=mixed
KillSignal=SIGTERM

# Security
NoNewPrivileges=true
ProtectHome=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mailngine-api

[Install]
WantedBy=multi-user.target
